# example_files/crontab
#crontab -e
0 * * * * /home/operator/Scripts/update-server.sh
