# .OnboardComputerSystem/config-files/crontab
@reboot   /home/laserwolf/.OnboardComputerSystem/lcd-boot.py
0 * * * * /home/laserwolf/Scripts/webcam.sh && /home/laserwolf/Scripts/update-server.sh
